# eap/azure/sr/sr.service
[Unit]
Description=RHPAM Smart Router service

[Service]
ExecStart=/bin/bash /opt/custom-config/sr-service.sh

[Install]
WantedBy=multi-user.target